# docker/python.ini
[program:python]
command= python /program/main.py